# wornwear-bot.service
# ── systemd service unit for Worn Wear Bot ────────────────────────────────────
# This keeps the bot running forever on a Linux VPS.
# If it crashes, systemd restarts it automatically.
#
# SETUP:
#   1. Copy this file to /etc/systemd/system/wornwear-bot.service
#   2. Edit the User, WorkingDirectory, and ExecStart paths below
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable wornwear-bot   # start on boot
#   5. sudo systemctl start wornwear-bot    # start now
#
# USEFUL COMMANDS:
#   sudo systemctl status wornwear-bot      # check if running
#   sudo journalctl -u wornwear-bot -f      # live logs
#   sudo systemctl stop wornwear-bot        # stop it
#   sudo systemctl restart wornwear-bot     # restart it

[Unit]
Description=Worn Wear Monitor Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/wornwear-bot
ExecStart=/home/ubuntu/wornwear-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
